# roles/lb/templates/smuggler-lb.service.j2
[Unit]
Description=Smuggler nftables rules for load balancing
After=network.target nftables.service
Requires=nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStartPre=/bin/bash -c 'nft list table ip lb >/dev/null 2>&1 && nft flush table ip lb || true'
ExecStart=/usr/sbin/nft -f /etc/nftables.d/smuggler-lb.conf
ExecStop=/bin/bash -c 'nft list table ip lb >/dev/null 2>&1 && nft flush table ip lb || true'
ExecStopPost=/bin/bash -c 'nft list table ip lb >/dev/null 2>&1 && nft delete table ip lb || true'

[Install]
WantedBy=multi-user.target
