# roles/tc/templates/ratelimit.service.j2
[Unit]
Description=Configure DNSTT ratelimit
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /opt/dnstt-ratelimit.sh

[Install]
WantedBy=multi-user.target
